# crates/portunus-standalone/contrib/portunus-standalone.service
[Unit]
Description=Portunus standalone TCP/UDP forwarder
Documentation=https://github.com/ZingerLittleBee/Portunus
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/portunus-standalone --config /etc/portunus/standalone.toml
# No config reload — modifying the TOML requires `systemctl restart portunus-standalone`.
# Omitting ExecReload makes `systemctl reload` correctly report "not applicable".
Restart=on-failure
RestartSec=2
StandardOutput=journal
StandardError=journal

# Resource limits — main runtime warns when nofile < 4096.
LimitNOFILE=65535

# User & filesystem hardening
User=portunus
Group=portunus
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
NoNewPrivileges=true
RuntimeDirectory=portunus
RuntimeDirectoryMode=0755
ReadOnlyPaths=/etc/portunus
ReadWritePaths=/run/portunus

# Network capabilities — allows binding privileged ports (<1024)
# without running as root. Required for ports such as 22, 53, 80, 443.
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
