# deploy/systemd/portunus-server.service
[Unit]
Description=Portunus control plane
Documentation=https://github.com/ZingerLittleBee/Portunus
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
User=portunus-server
Group=portunus-server
ExecStart=/usr/local/bin/portunus-server --data-dir /var/lib/portunus serve
Restart=on-failure
RestartSec=5

# Drain in-flight forwarded connections before SIGKILL. Keep this larger
# than the server.toml `shutdown_drain_timeout_secs` (default 30 s) so the
# server's own drain logic gets to run.
TimeoutStopSec=60

# File-descriptor headroom for port-range rules (002-port-range-forward).
# The server itself does NOT bind range ports (the client does), so the
# server only needs fds for: gRPC streams (1 per connected client) +
# operator HTTP + metrics + persistence + headroom. 4096 is comfortable
# for ~hundreds of clients. Bumped from the soft default (1024) so a
# rule.list rendering doesn't run into a transient fd limit.
LimitNOFILE=4096

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictAddressFamilies=AF_INET AF_INET6
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM

# Writable: SQLite state, generated TLS material, optional server.toml.
ReadWritePaths=/var/lib/portunus

[Install]
WantedBy=multi-user.target
